# sla-dashboard.service
[Unit]
Description=SLA Monitoring Dashboard
After=network.target

[Service]
Type=simple
User=Nfsscid
WorkingDirectory=/data/sla-dashboard
ExecStart=/data/sla-dashboard/venv/bin/python /data/sla-dashboard/main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Keep uploaded file in memory across crashes is not possible,
# but the app will still serve the dashboard on restart.
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
